# templates/default/razor.erb
#! /bin/bash

### BEGIN INIT INFO
# Provides:          razor
# Required-Start:    $network $mongodb $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Razor 
# Description:       Razor provisioning service
### END INIT INFO

PATH=<%= @directory %>/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
export PATH

DAEMON=<%= @directory %>/bin/razor_daemon.rb

test -x $DAEMON || exit 0

# We can just delegate all our operations to the Ruby wrapper script.
exec $DAEMON "$@"
